# supervisor/app.conf
[program:app]
directory=/opt/dash_app
command=python3 app.py
